# sync/com.iselodyhereyet.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.iselodyhereyet.sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/zach/Code/iselodyhereyet/sync/sync.sh</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>SYNC_AUTH_TOKEN</key>
		<string>REPLACE_ME</string>
		<key>SYNC_URL</key>
		<string>https://iselodyhereyet-sms.zachjustice.workers.dev/sync</string>
	</dict>

	<key>StartInterval</key>
	<integer>300</integer>

	<key>StandardOutPath</key>
	<string>/tmp/iselodyhereyet-sync.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/iselodyhereyet-sync.log</string>

	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
